# roles/atproto_handle_check/templates/atproto_handle_check.timer.j2
[Unit]
Description=atproto/Bluesky handle DNS verification timer

[Timer]
OnCalendar={{ atproto_handle_check_oncalendar }}
Persistent=true

[Install]
WantedBy=timers.target
