# roles/caddy/templates/caddy_certs.path.j2
[Unit]
Description=Caddy watch for new certs
After=network-online.target
StartLimitInterval=0
StartLimitIntervalSec=0

[Path]
PathModified={{ caddy_certs_modified_path }}
Unit=caddy_certs.service
